# deploy/linux/CLIProxyAPIManager.service
[Unit]
Description=CLIProxyAPIManager
Documentation=https://github.com/a389216533/CLIProxyAPIManager
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# The README installation command replaces this placeholder with the package directory.
WorkingDirectory=__CLIPROXYAPIMANAGER_DIR__
ExecStart=/bin/sh -c 'exec ./CLIProxyAPIManager'
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
